# systemd/vendorhub-producer-b.service
[Unit]
Description=VendorHub Producer B (CSV -> staging + RMQ)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/vendor-sync/app
EnvironmentFile=/opt/vendor-sync/.env
ExecStart=/opt/vendor-sync/app/venv/bin/python -m vendorhub.cli.run_producer_b
User=root
Group=root
