# systemd/vendorhub-cleanup.timer
[Unit]
Description=Run VendorHub Cleanup daily at 02:00 UTC

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true
Unit=vendorhub-cleanup.service

[Install]
WantedBy=timers.target
